# docker/s6-rc.d/dashboard/run
#!/bin/sh
# shellcheck shell=sh
# Dashboard service. Always declared so s6 has a supervised slot; if
# HERMES_DASHBOARD isn't truthy the run script exits cleanly and the
# companion finish script returns 125 (s6's "permanent failure, do
# not restart" marker), so s6-svstat reports the slot as down. See
# also docker/s6-rc.d/dashboard/finish.
#
# NOTE: This run script no longer uses `#!/command/with-contenv sh`.
# s6-overlay's with-contenv has been observed to drop env vars
# (notably ADMIN_PASSWORD, API_SERVER_KEY) when running on Railway,
# leaving supervised services without container-injected secrets even
# though the dashboard CMD inherits them fine. We source /proc/1/environ
# instead via the helper below — see docker/lib/source-container-env.sh.

# Re-import the full container env from PID 1 before anything reads it.
# shellcheck disable=SC1091
. /dev/stdin <<EOF
$(python3 - <<'PY'
import shlex
try:
    with open('/proc/1/environ', 'rb') as f:
        blob = f.read()
except OSError:
    blob = b''
for entry in blob.split(b'\x00'):
    if not entry or b'=' not in entry:
        continue
    name, _, value = entry.partition(b'=')
    try:
        name_s = name.decode('utf-8')
        value_s = value.decode('utf-8')
    except UnicodeDecodeError:
        continue
    if not name_s or not (name_s[0].isalpha() or name_s[0] == '_'):
        continue
    if not all(c.isalnum() or c == '_' for c in name_s):
        continue
    print(f"export {name_s}={shlex.quote(value_s)}")
PY
)
EOF

case "${HERMES_DASHBOARD:-}" in
    1|true|TRUE|True|yes|YES|Yes) ;;
    *)
        # Exit 0; the finish script will exit 125 → s6-supervise won't
        # restart us and the slot reports down. Using a clean exit
        # (rather than `exec sleep infinity`) means s6-svstat reflects
        # reality: when HERMES_DASHBOARD is unset, the service is NOT
        # running, just supervised-with-permanent-failure. See PR
        # #30136 review item I3.
        exit 0
        ;;
esac

cd /opt/data
# shellcheck disable=SC1091
. /opt/hermes/.venv/bin/activate

dash_host="${HERMES_DASHBOARD_HOST:-0.0.0.0}"
dash_port="${HERMES_DASHBOARD_PORT:-9119}"

# Binding to anything other than localhost requires --insecure — the
# dashboard refuses otherwise because it exposes API keys. Inside a
# container this is the expected deployment.
insecure=""
case "$dash_host" in
    127.0.0.1|localhost) ;;
    *) insecure="--insecure" ;;
esac

# Pre-bind probe: if another process (typically the main CMD's
# run_web_server.py) is already serving on this port, exit cleanly and
# let the finish script mark this slot as permanently down. Without
# this, s6 would restart us in a tight loop on EADDRINUSE.
if ! python3 -c "
import socket, sys
s = socket.socket(socket.AF_INET, socket.SOCK_STREAM)
try:
    s.bind(('0.0.0.0', int('${dash_port}')))
except OSError:
    sys.exit(1)
finally:
    s.close()
"; then
    echo "  dashboard service: port ${dash_port} already bound — going idle." >&2
    # Sleep forever so s6 keeps the slot 'up' (avoids restart spam) but
    # we never try to bind. The finish script handles termination.
    exec sleep infinity
fi

# shellcheck disable=SC2086  # word-splitting of $insecure is intentional
exec s6-setuidgid hermes hermes dashboard \
    --host "$dash_host" --port "$dash_port" --no-open $insecure
